# tags/RELEASE/rc-scripts-0.4/init.d/hostname
#!/sbin/runscript
# Copyright © 2005  TUBITAK/UEKAE
# Licensed under the GNU General Public License, version 2.
# See the file http://www.gnu.org/copyleft/gpl.txt.
#
# Original work belongs Gentoo Linux

depend() {
	need checkroot
}

start() {
	local myhost=$(/bin/hostname 2>/dev/null)
	local retval=0

	# If the hostname is already set via the kernel, and /etc/hostname 
	# isn't setup, then we shouldn't go reseting the configuration #38172.
	if [[ -z ${myhost} ]] || [[ ${myhost} == "(none)" ]] ; then
		myhost="localhost"
	fi

	if [[ -f /etc/hostname ]] ; then
		ewarn "You should stop using /etc/hostname and use /etc/conf.d/hostname"
		myhost=$(</etc/hostname)
	else
		myhost=${HOSTNAME}
	fi

	ebegin "Setting hostname"
	/bin/hostname "${myhost}"
	retval=$?
	eend ${retval} "Failed to set the hostname"

	if [[ ${retval} -eq 0 ]] ; then
		# setup $HOSTNAME, ignore errors in case /etc is readonly.
		echo "HOSTNAME=\"${myhost}\"" 2>/dev/null > /etc/env.d/01hostname
	fi

	return ${retval}
}

# vim:ts=4
