# tags/RELEASE/rc-scripts-0.4/init.d/local
#!/sbin/runscript
# Copyright © 2005  TUBITAK/UEKAE
# Licensed under the GNU General Public License, version 2.
# See the file http://www.gnu.org/copyleft/gpl.txt.
#
# Original work belongs Gentoo Linux

depend() {
	after *
}

start() {
	ebegin "Starting local"

	# Add any misc programs that should be started
	# to /etc/conf.d/local.start
	if [[ -e /etc/conf.d/local.start ]] ; then
		source /etc/conf.d/local.start
	fi

	eend $? "Failed to start local"
}

stop() {
	ebegin "Stopping local"

	# Add any misc programs that should be stopped
	# to /etc/conf.d/local.stop
	if [[ -e /etc/conf.d/local.stop ]] ; then
		source /etc/conf.d/local.stop
	fi

	eend $? "Failed to stop local"
}


# vim:ts=4
