# tags/RELEASE/rc-scripts-0.4/init.d/rmnologin
#!/sbin/runscript
# Copyright © 2005  TUBITAK/UEKAE
# Licensed under the GNU General Public License, version 2.
# See the file http://www.gnu.org/copyleft/gpl.txt.
#
# Original work belongs Gentoo Linux

depend() {
	need localmount
}

start() {
	if [[ -f /etc/nologin.boot ]] ; then
		rm -f /etc/nologin /etc/nologin.boot &> /dev/null
	fi
}


# vim:ts=4
