# tags/RELEASE/rc-scripts-0.3/init.d/numlock
#!/sbin/runscript
# Copyright © 2005  TUBITAK/UEKAE
# Licensed under the GNU General Public License, version 2.
# See the file http://www.gnu.org/copyleft/gpl.txt.
#
# Original work belongs Gentoo Linux

depend() {
	need localmount
}

start() {
	ebegin "Enabling numlock on ttys"
	local dev
	[[ -d /dev/vc ]] \
		&& dev=/dev/vc/ \
		|| dev=/dev/tty
	for tty in $(seq 1 "${RC_TTY_NUMBER}") ; do
		setleds -D +num < ${dev}${tty} &> /dev/null
	done
	eend $? "Failed to enable numlock"
}

stop() {
	ebegin "Disabling numlock on ttys"
	local dev
	[[ -d /dev/vc ]] \
		&& dev=/dev/vc/ \
		|| dev=/dev/tty
	for tty in $(seq 1 "${RC_TTY_NUMBER}") ; do
		setleds -D -num < ${dev}${tty} &> /dev/null
	done
	eend $? "Failed to disable numlock"
}


# vim:ts=4
